# config/crontab
# Run the script every day at 23:59 PM
59 23 * * * /usr/bin/python3 /app/yyj_scraper.py